# myservice.service
[Unit]
Description=My Service

[Service]
Type=simple
ExecStart=/home/eren-linuxmint/Masaüstü/BSM/script.sh

[Install]
WantedBy=multi-user.targe
